# root/etc/services.d/01-mount-media/finish
#!/usr/bin/with-contenv sh

umount_command="${UMOUNT_MEDIA}"

echo "Executing $umount_command"
exec $umount_command

echo "unmounted media"