# root/etc/services.d/01-mount-media/run
#!/usr/bin/with-contenv sh

mount_command="${MOUNT_MEDIA}"

echo "Started ${mount_command}"
exec $mount_command