# root/etc/services.d/01-mount-rclone/finish
#!/usr/bin/with-contenv sh

umount_command="fusermount $UnmountCommands /drive-media"

echo "Executing $umount_command"
exec $umount_command

echo "mounted"